# infra/com.aria.octant-chrome.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<!--
		FALLBACK ONLY — Use this if headless Chromium + MetaMask fails.

		This LaunchAgent starts Chrome with remote debugging enabled in the GUI
		session context (Aqua namespace). Requires auto-login on Mac Mini.

		Install:
		  cp infra/com.aria.octant-chrome.plist ~/Library/LaunchAgents/
		  launchctl load ~/Library/LaunchAgents/com.aria.octant-chrome.plist

		Verify:
		  curl http://localhost:9222/json

		Then configure Playwright to connect via CDP:
		  browserType.connectOverCDP('http://localhost:9222')
	-->
	<key>Label</key>
	<string>com.aria.octant-chrome</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Google Chrome.app/Contents/MacOS/Google Chrome</string>
		<string>--remote-debugging-port=9222</string>
		<string>--no-first-run</string>
		<string>--no-default-browser-check</string>
		<string>--disable-default-apps</string>
		<string>--user-data-dir=/tmp/octant-chrome-profile</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/tmp/octant-chrome-stdout.log</string>

	<key>StandardErrorPath</key>
	<string>/tmp/octant-chrome-stderr.log</string>
</dict>
</plist>
